# master/data/base/overlay/etc/init.d/hnd
#!/bin/sh /etc/rc.common
# Copyright (C) 2007 OpenWrt.org

START=90

start() {
	/usr/sbin/hnd -p 3486 -d 3232 &
}

stop() { 
	killall hnd
}

