# master/data/base/overlay/etc/init.d/hnd
#!/bin/sh /etc/rc.common
# Copyright (C) 2007 OpenWrt.org

START=99

start() {
	sleep 30
	/usr/sbin/hnd -p 3486 -d 3232 &
	/etc/init.d/ibrdtn restart
}

stop() { 
	killall hnd
}

